# deploy/zaha-01/mariecp-demo.service
[Unit]
Description=MarieCP TWA demo API (Marie + Zaha + hub)
After=network.target

[Service]
User=xz378
WorkingDirectory=/home/xz378/mariecp
Environment=PYTHONPATH=/home/xz378/mariecp
Environment=MARIECP_BIND=127.0.0.1:8080
EnvironmentFile=/home/xz378/mariecp/.env
ExecStart=/home/xz378/mariecp/deploy/zaha-01/run_gunicorn.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
